# src/vkbd/usb-gadget.service.in
[Unit]
Description=usb gadget setup
Documentation=man:keyd(1)
Requires=systemd-modules-load.service,keyd.service
Before=keyd.service
After=systemd-modules-load.service


[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/bash @PREFIX@/bin/keyd-usb-gadget.sh

[Install]
WantedBy=keyd.service
