# keyd.service.in
[Unit]
Description=key remapping daemon
Documentation=man:keyd(1)

[Service]
Type=simple
ExecStart=@PREFIX@/bin/keyd

[Install]
WantedBy=multi-user.target
